# systemd.service
[Unit]
Description=YourApp
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/your_app
ExecStart=/home/deploy/your_app/bin/your_app foreground
ExecStop=/home/deploy/your_app/bin/your_app stop
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/deploy/your_app/config/env.prod
SyslogIdentifier=yourapp
RemainAfterExit=no

[Install]
WantedBy=multi-user.target
